# supervisord.single.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:surrealdb]
command=surreal start --log trace --user root --pass root rocksdb:/mydata/mydatabase.db
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=5
autostart=true
startsecs=5

[program:api]
# 添加 keep-alive 和超時設置以支持長時間運行的請求
# --timeout-keep-alive 300: 保持連接活躍 300 秒（5 分鐘）
# --timeout-graceful-shutdown 30: 優雅關閉超時 30 秒
command=uv run uvicorn api.main:app --host 0.0.0.0 --port 5055 --timeout-keep-alive 300 --timeout-graceful-shutdown 30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=10
autostart=true
startsecs=3

[program:worker]
command=uv run surreal-commands-worker --import-modules commands
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=20
autostart=true
startsecs=3

[program:frontend]
command=bash -c "sleep 5 && npm run start"
directory=/app/frontend
environment=NODE_ENV="production",PORT="8188"
passenv=API_URL,NEXT_PUBLIC_API_URL,INTERNAL_API_URL
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=30
autostart=true
startsecs=5